# systemd/psql-backup.service
[Unit]
Description=weekley pg_dumpall

[Service]
Type=simple
ExecStart=/hoge/psql-backup.sh
User=hogehoge